# dist/paralet.service
[Unit]
Description=Parasaus paralet
After=network.target

[Service]
Type=simple
ExecStart=-/usr/bin/paralet --etcd-servers=http://127.0.0.1:2379
#ExecStart=-/usr/bin/paraletsh
ExecStop=/bin/kill -WINCH ${MAINPID}

Restart=always
RestartSec=10s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target